# supervisord.conf
[program:laravel-worker]
command=php artisan queue:work --daemon
directory=/var/www/html
autostart=true
autorestart=true
stderr_logfile=/var/www/html/storage/logs/supervisor.log
stdout_logfile=/var/www/html/storage/logs/supervisor.log


[program:laravel-scheduler]
command=php /var/www/html/artisan schedule:run
directory=/var/www/html
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stderr_logfile=/var/www/html/storage/logs/scheduler-error.log
autostart=true
autorestart=true
startretries=3
user=root


